# scripts/systemd/milton-gateway.service.template
[Unit]
Description=Milton Chat Gateway (OpenAI-compatible API)
Documentation=https://github.com/cole-hanan/milton
After=network.target milton-api.service
Wants=milton-api.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory={{MILTON_REPO_DIR}}
EnvironmentFile={{MILTON_REPO_DIR}}/scripts/milton.env

# Python executable will be substituted by milton_up.sh
ExecStart={{PYTHON_EXECUTABLE}} {{MILTON_REPO_DIR}}/scripts/start_chat_gateway.py

# Logging
StandardOutput=append:{{MILTON_STATE_DIR}}/logs/milton-gateway.log
StandardError=append:{{MILTON_STATE_DIR}}/logs/milton-gateway.log

# Restart policy
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
