# scripts/systemd/milton-morning-briefing.service
[Unit]
Description=Milton Morning Briefing
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/milton
Environment="PATH=%h/miniconda3/envs/milton/bin:/usr/local/bin:/usr/bin"
ExecStart=/bin/bash -c '%h/miniconda3/envs/milton/bin/python3 %h/milton/scripts/enhanced_morning_briefing.py && %h/miniconda3/envs/milton/bin/python3 %h/milton/scripts/send_briefing_to_phone.py --method ntfy --briefing morning'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
