# systemd/milton-phone-listener.service
[Unit]
Description=Milton Phone Question Listener (ntfy)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/home/cole-hanan/milton
Environment="STATE_DIR=%h/.local/state/milton"
EnvironmentFile=/home/cole-hanan/milton/.env
Environment="PATH=/home/cole-hanan/miniconda3/envs/milton/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/cole-hanan/milton/scripts/ask_from_phone.py --listen
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
