# systemd/milton-reminders.service
[Unit]
Description=Milton Reminders Scheduler
Documentation=https://github.com/colehanan1/milton/blob/main/docs/reminders.md
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/milton
Environment="PATH=%h/miniconda3/envs/milton/bin:/usr/local/bin:/usr/bin"
Environment="STATE_DIR=%h/.local/state/milton"
Environment="NTFY_BASE_URL=https://ntfy.sh"
Environment="NTFY_TOPIC=milton-reminders-cole"
Environment="TZ=America/New_York"
ExecStart=%h/miniconda3/envs/milton/bin/milton-reminders run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
